# containers/ureport-celery/conf/s6-rc.d/celery-beat/run
#!/command/with-contenv sh

cd /var/www

if test "$DEBUG" = true; then
    LOG_LEVEL="debug"
else
    LOG_LEVEL="warning"
fi

celery -A ureport beat \
    --loglevel $LOG_LEVEL
