# rapidpro/docker/s6-rc.d/nginx/run
#!/command/with-contenv sh

set -e

envsubst '$MAILROOM_URL $COURIER_URL' </etc/nginx/templates/default >/etc/nginx/sites-available/default
nginx -g "daemon off;"
